# .deploy/php/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:php-fpm]
command=php-fpm

[program:horizon]
command=php artisan horizon
user=www-data
